# io.opencaesar.rosetta.docker/app/conf.d/rosetta.conf
[program:rosetta]
command=/app/Rosetta/Rosetta -data /data/workspace
autorestart=false
